# target/linux/ipq40xx/base-files/etc/init.d/mt7615-rmmod
#!/bin/sh /etc/rc.common

STOP=52

stop() {
	rmmod mt7615e > /dev/null 2>&1
}
